# debian/lingmo-sdk-location.postinst
#!/bin/bash

systemctl enable dbus_location.service
systemctl restart dbus_location.service